# web_console_v2/api/cmds/supervisord_dev.conf
[supervisord]
pidfile=/tmp/supervisord.pid
loglevel=debug

; The section sets up an HTTP server that listens on "file", which can be used
; to control the daemon
[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; This section connects to the HTTP server to control the HTTP server
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:restful_api]
command=./web_console_v2/api/cmds/flask_cli_bin run --eager-loading --port=1991 --host=0.0.0.0
environment=FLASK_APP=web_console_v2.api.server:app
redirect_stderr=true

[program:rpc]
command=./web_console_v2/api/entrypoint_bin -- start-rpc
redirect_stderr=true

[program:composer]
command=./web_console_v2/api/entrypoint_bin -- start-composer
redirect_stderr=true
