# web_console_v2/api/cmds/supervisord.conf
[supervisord]
pidfile=/run/supervisord.pid
loglevel=debug

; The section sets up an HTTP server that listens on "file", which can be used
; to control the daemon
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; This section connects to the HTTP server to control the HTTP server
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:restful_api]
command=web_console_v2/api/cmds/gunicorn_cli_bin --conf web_console_v2/api/gunicorn_config.py server:app
redirect_stderr=true

[program:rpc]
command=web_console_v2/api/entrypoint_bin start-rpc
redirect_stderr=true

[program:composer]
command=web_console_v2/api/entrypoint_bin start-composer
redirect_stderr=true
